# ai.openclaw.eth-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.eth-watcher</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>__RUN_SCRIPT__</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__WORKING_DIR__</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>/tmp/eth-invest-agent-daemon.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/eth-invest-agent-daemon.err.log</string>
</dict>
</plist>
